# pkg/deb/prerm
#!/bin/bash
set -e

case "$1" in
    remove)
        echo "natman-go: Stopping and disabling service..."
        systemctl stop natman 2>/dev/null || true
        systemctl disable natman 2>/dev/null || true
        ;;
    upgrade)
        echo "natman-go: Stopping service for upgrade..."
        systemctl stop natman 2>/dev/null || true
        ;;
esac